# deploy/bybit-pump-scanner.service
[Unit]
Description=Bybit Pump Exhaustion Scanner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/bybit-scanner
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/bybit-scanner/main_pump.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
